# examples/traffic-count.service
[Unit]
Description=mac-traffic-ebpf traffic counting service
After=network.target

[Service]
Type=simple
# Run as root; eBPF TC attachment requires CAP_BPF + CAP_NET_ADMIN
# AmbientCapabilities grants these without full root privilege escalation
AmbientCapabilities=CAP_BPF CAP_NET_ADMIN
CapabilityBoundingSet=CAP_BPF CAP_NET_ADMIN
NoNewPrivileges=false
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536

# Hardening
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/traffic-count
ReadWritePaths=/var/lib/traffic-count

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=traffic-count

[Install]
WantedBy=multi-user.target
